# gcs-qgc/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1440x900x24
autorestart=true
priority=10

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":0"
autorestart=true
priority=20

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -nopw -forever -shared -rfbport 5900
autorestart=true
priority=30

[program:novnc]
# nginx 가 8080 reverse proxy 하므로 websockify 는 내부 8081 로 옮긴다.
command=/usr/bin/websockify --web=/usr/share/novnc/ 8081 localhost:5900
autorestart=true
priority=40

[program:nginx]
# UAVGcsAccess_CL NDJSON access log → 컨테이너 stdout (PID 1 = supervisord)
# → fluent-bit kube tail → DCR. nginx daemon off 필수.
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=45
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:qgc]
command=/opt/qgc/AppRun
user=qgc
environment=DISPLAY=":0",QT_X11_NO_MITSHM="1",HOME="/home/qgc"
autorestart=true
priority=50
